# lovense-mcp.service
[Unit]
Description=Lovense MCP Server for Nate's Consciousness Substrate
Documentation=https://github.com/clotho2/nate_api_substrate
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aicara/wolfe-lovense

# Environment variables
Environment="LOVENSE_DEVELOPER_TOKEN=YOUR_TOKEN_HERE"
Environment="LOVENSE_CALLBACK_URL=http://94.130.203.232:8000/lovense/callback"
Environment="PYTHONUNBUFFERED=1"

# Start command
ExecStart=/opt/aicara/wolfe-lovense/venv/bin/python /opt/aicara/wolfe-lovense/Enhanced_Lovense_MCP.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aicara/wolfe-lovense

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lovense-mcp

[Install]
WantedBy=multi-user.target
